# src/promptflow/promptflow/_sdk/data/docker_csharp/runit/promptflow-serve/run
#! /bin/bash

echo "start promptflow serving"
cd /flow
dotnet Promptflow.DotnetService.dll -p "8080" -y "flow.dag.yaml" -a "." -c "" -l "" -s